# mobsboi-parser.service
[Unit]
Description=МобСбой — парсер Downdetector
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/mobsboi
ExecStart=/var/www/mobsboi/venv/bin/python parser.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/mobsboi-parser.log
StandardError=append:/var/log/mobsboi-parser.log

[Install]
WantedBy=multi-user.target
