# mobsboi-server.service
[Unit]
Description=МобСбой — веб-сервер FastAPI
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/mobsboi
ExecStart=/var/www/mobsboi/venv/bin/uvicorn server:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=append:/var/log/mobsboi-server.log
StandardError=append:/var/log/mobsboi-server.log

[Install]
WantedBy=multi-user.target
